# config_files/systemd/hostapd.service.d/requiresdevice.conf
[Unit]
BindsTo=sys-subsystem-net-devices-wlx0002728c289a.device
After=sys-subsystem-net-devices-wlx0002728c289a.device


